# uberspaceconfig/crontab/descriptors_simple_cronjob.ini
# Descriptor Computation Cronjob
# Runs every 5 minutes to compute missing descriptors for components
*/10 * * * * /home/ddu/csc/venv/bin/python3.9 /home/ddu/csc/backend/main_descriptors_simple.py >> /home/ddu/csc/backend/logs/descriptors_simple_cronjob.log 2>&1