# uberspaceconfig/crontab/descriptorcomputation_cronjob.ini
# Descriptor Computation Cronjob
# Runs every 5 minutes to compute missing descriptors for components
*/5 * * * * cd ~/csc/src/backend && ~/miniconda3/envs/csc/bin/python main_computedescriptors.py >> ~/logs/descriptorcomputation.log 2>&1