# uberspaceconfig/crontab/geometrymaintenance_cronjob.ini
MAILTO="eschenbach@dg.tu-darmstadt.de"
0 3 * * * /home/ddu/csc/venv/bin/python3.9 /home/ddu/csc/backend/geometrymaintenance.py >> /home/ddu/csc/backend/logs/geometrymaintenance_cronjob.log 2>&1