# uberspaceconfig/crontab/ghxml_sync_cronjob.ini
# GH XML Sync Cronjob to mirror pasteable XML files from private repo
# Runs every 30 minutes
MAILTO="eschenbach@dg.tu-darmstadt.de"
*/30 * * * * /bin/bash /home/ddu/csc/backend/ghxml_sync.sh >> /home/ddu/csc/backend/logs/ghxml_sync.log 2>&1
