# uberspaceconfig/crontab/usermaintenance_cronjob.ini
MAILTO="eschenbach@dg.tu-darmstadt.de"
0 2 * * * /home/ddu/csc/venv/bin/python3.9 /home/ddu/csc/backend/usermaintenance.py >> /home/ddu/csc/backend/logs/usermaintenance_cronjob.log 2>&1